# triggers/Library/LaunchDaemons/com.github.patchoo-trigger-every120.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.github.patchoo-trigger-every120</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/local/bin/jamf</string>
		<string>policy</string>
		<string>-trigger</string>
		<string>every120</string>
		<string>-randomDelaySeconds</string>
		<string>300</string>
	</array>
	<key>StartInterval</key>
	<integer>7200</integer>
	<key>UserName</key>
	<string>root</string>
</dict>
</plist>
